# experts/orchestrator_1_0_prompt_evaluation/orchestrator_api.service
[Unit]
Description=Orchestrator API Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/FrizzlesRubric/experts/orchestrator_1_0_prompt_evaluation
ExecStart=/bin/bash -c 'source /home/ubuntu/FrizzlesRubric/experts/orchestrator_1_0_prompt_evaluation/venv/bin/activate && uvicorn orchestrator_api:app --host 0.0.0.0 --port 8010'
Restart=on-failure

[Install]
WantedBy=multi-user.target 